# packages/system/debian/etc/systemd/system/ssh.service.d/local.conf
[Unit]
Wants=sn-sshdgenkeys.service
After=network.target auditd.service sn-sshdgenkeys.service
