# packages/nftables/debian/etc/cron.d/nftables-drop-brute
*/2 * * * * root /usr/share/solarnode/bin/nftables-drop-brute.sh >/dev/null 2>&1
