# packages/pi-sixfab-3g-shield-vodafone/debian/lib/systemd/system/sn-sixfab-3g-pppd.service
[Unit]
Description=PPPd managed connection

[Service]
ExecStart=/usr/sbin/pppd call provider nodetach nolog
Restart=on-failure
RestartSec=30
