# packages/pi-sixfab-3g-shield-vodafone/debian/lib/systemd/system/sn-sixfab-3g-reconnect.service
[Unit]
Description=Mobile network auto reconnect

[Service]
Type=oneshot
ExecStart=/usr/share/solarnode/bin/sixfab-3g-reconnect.sh
KillMode=none
