# packages/iptables/debian/sn-iptables.postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
	systemctl daemon-reload >/dev/null || true

	if ! systemctl is-enabled sn-iptables.service >/dev/null; then
		echo 'Enabling sn-iptables.service...'
		systemctl enable sn-iptables.service || true
	
		echo 'iptables will not be enforced until a reboot, or if you run `systemctl start sn-iptables`.'
	fi
fi
