# solarssh/debian/lib/systemd/system/solarssh-cleaner.timer
[Unit]
Description=Close discarded SolarSSH sessions daily.

[Timer]
OnCalendar=hourly
RandomizedDelaySec=120

[Install]
WantedBy=timers.target

