# wifi/debian/lib/systemd/system/sn-wifi-conf@.path
[Unit]
Description=WPA supplicant configuration file (interface-specific version)
Requires=sys-subsystem-net-devices-%i.device
After=sys-subsystem-net-devices-%i.device
Before=network.target
Wants=network.target

[Path]
Unit=sn-wifi-conf@%i.service
PathExists=/etc/wpa_supplicant/wpa_supplicant-%i.conf
PathChanged=/etc/wpa_supplicant/wpa_supplicant-%i.conf

[Install]
Alias=multi-user.target.wants/sn-wifi-conf@%i.path

