# packages/iptables/debian/etc/cron.d/drop-brute
*/2 * * * * root /usr/share/solarnode/bin/drop-brute.sh >/dev/null 2>&1
