# pi-sixfab-3g-shield-vodafone/debian/lib/systemd/system/sn-sixfab-3g-reconnect.timer
[Unit]
Description=Mobile network auto reconnect periodic task

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
RandomizedDelaySec=10

[Install]
WantedBy=timers.target
