# packages/system/debian/lib/systemd/system/sn-expandfs.service
[Unit]
Description=SolarNode Filesystem Expand
After=sysinit.target local-fs.target
ConditionPathExists=/boot/sn-expandfs

[Service]
Type=oneshot
ExecStart=/usr/share/solarnode/bin/expandfs.sh -v -o /var/local/sn-expandfs.saved
ExecStartPost=/bin/rm -f /boot/sn-expandfs

[Install]
WantedBy=basic.target
