# motion/debian/lib/systemd/system/sn-motion-cleaner.timer
[Unit]
Description=Clean old motion images daily.

[Timer]
OnCalendar=daily
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
