# motion/debian/lib/systemd/system/sn-motion-cleaner.service
[Unit]
Description=Clean out old motion media files

[Service]
Type=oneshot
ExecStart=/usr/share/solarnode/bin/clean-motion-oldmedia.sh
