# system/debian/lib/systemd/system/sn-reset-app.service
[Unit]
Description=Reset SolarNode (application only)

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/share/solarnode/bin/reset.sh -a'
