# sys/freebsd/etc/periodic/weekly/999.sn-clean-logs-postgres
#!/bin/sh
/usr/local/share/sn/bin/clean-logs.sh -d 35 -v $@ /var/log/postgres

