# packages/pi/debian/lib/udev/rules.d/15-i2c-solarnode-pi.rules
SUBSYSTEM=="i2c|spi", ENV{MODALIAS}=="?*", ENV{OF_NAME}=="?*", ENV{OF_COMPATIBLE_0}=="?*", RUN+="/usr/share/solarnode/bin/i2cprobe.sh"
