# packages/pi-sixfab-3g-shield-vodafone/debian/sn-pi-sixfab-3g-shield-vodafone.prerm
#!/bin/sh -e

service_down () {
	name="$1"
	if systemctl is-enabled "$name" >/dev/null; then
		echo "Disabling $name..."
		systemctl disable --now "$name" || true
	fi
}

if [ "$1" = "remove" ]; then
	service_down sn-sixfab-3g-reconnect.timer
	service_down sn-sixfab-3g-pppd
fi
